# .github/dib-elements/custom-deploy/overlay/etc/systemd/system/ironic-python-agent.service.d/override.conf
[Service]
# Uncomment to enable debug logging
#ExecStart=
#ExecStart=/usr/local/bin/ironic-python-agent --debug --config-dir /etc/ironic-python-agent.d/
